# com.briansmith.flickrupload.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.briansmith.flickrupload</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/bsmi067/Claude/Projects/flickr_update/run_upload.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Minute</key>
        <integer>17</integer>
    </dict>
    <key>StandardErrorPath</key>
    <string>/Users/bsmi067/Claude/Projects/flickr_update/launchd_error.log</string>
</dict>
</plist>
